# db/crontab
59 23 * * * /home/gsrinivasan/Code/Fraud-net/db/venv/bin/python3 /home/gsrinivasan/Code/Fraud-net/db/deletion.py
